# etc/hotplug.d/iface/35-dibbler
#!/bin/sh
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && /etc/init.d/dibblerclient enabled && /etc/init.d/dibblerclient start
